# CRONTAB/crontab
# Ежедневное резервное копирование домашней директории в 12:00
0 12 * * * /usr/local/bin/backup_home.sh
# Для тестирования - каждые 10 минут
*/10 * * * * /usr/local/bin/backup_homefolder.sh